# i3/.config/systemd/user/i3.target
[Unit]
Description=Window Manager i3
Requires=default.target
After=default.target
Wants=redshift.service
Wants=compton.service
Wants=dunst.service
