# systemd-user/.config/systemd/user/compton.service
[Unit]
Description=X Window Compositor

[Service]
Type=simple
ExecStart=/usr/bin/compton --dbus
KillMode=mixed
KillSignal=SIGINT

