# systemd-user/.config/systemd/user/dunst.service
[Unit]
Description=Lightweight Notification Daemon

[Service]
Type=simple
ExecStart=/usr/bin/dunst
KillMode=mixed
KillSignal=SIGINT
